# servers/static/etc/init.d/russ
#!/bin/sh
#
# Setup for russ-based services.
#
### BEGIN INIT INFO
# Provides:          russ
# Required-Start:    $syslog $time
# Required-Stop:     $syslog $time
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Sets up for russ-based services
# Description:       Sets up for russ-based services
### END INIT INFO

RUN_DIR="/var/run/russ"

start() {
	# ensure working directories
	mkdir -m 0755 -p "${RUN_DIR}"
	mkdir -m 0755 -p "${RUN_DIR}/pids"
	mkdir -m 0755 -p "${RUN_DIR}/services"

	# start services
	ruservice . start
}

stop() {
	ruservice . stop
}

case "$1" in
restart)
	stop
	sleep 1
	start
	;;
start)
	start
	;;
status)
	ruservice . status
	;;
stop)
	stop
	;;
*)
	echo "Usage: /etc/init.d/russ start|stop|restart|status"
	exit 2
	;;
esac
exit 0
